# freebsd/serenitynow
#!/bin/sh 

# PROVIDE: serenitynow
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf to enable serenitynow:
#
# serenitynow_enable="YES"
#
#

. /etc/rc.subr 
name="serenitynow" 
rcvar=${name}_enable

pidfile="/var/run/${name}.pid"
serenitynow_env_file="/root/.${name}.env"
serenitynow_command="/usr/local/bin/${name}"

required_files="/usr/local/share/${name}.d/${name}.sqlite3 /usr/local/bin/${name} ${serenitynow_env_file}" 

sig_stop="TERM"

command="/usr/sbin/daemon"
command_args="--close-fds --restart --user ${name} --supervisor-pidfile ${pidfile} --output-file /var/log/${name} ${serenitynow_command}"

load_rc_config $name
run_rc_command "$1"
